# services/hyperion-router.service
[Unit]
Description=Hyperion Router - Telegram to Claude Code bridge
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/lobster
Environment=PATH=/root/.local/bin:/root/.cargo/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/root/lobster/config/config.env
EnvironmentFile=-/config.env
ExecStart=/root/lobster/.venv/bin/python /root/lobster/src/bot/lobster_bot.py
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyperion-router

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
