# services/hyperion-claude.service
[Unit]
Description=Hyperion Claude - Always-on Claude Code session
After=network.target hyperion-router.service
Wants=hyperion-router.service

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/root/lobster-workspace
Environment=PATH=/root/.claude/bin:/root/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/root
ExecStart=/usr/bin/tmux -L hyperion new-session -d -s hyperion -c /root/lobster-workspace /root/lobster/scripts/claude-wrapper.exp
ExecStop=/usr/bin/tmux -L hyperion kill-session -t hyperion
RemainAfterExit=yes
Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyperion-claude

[Install]
WantedBy=multi-user.target
