# services/lobster-claude.service.template
[Unit]
Description=Lobster Claude - Always-on Claude Code session
After=network.target lobster-router.service
Wants=lobster-router.service

[Service]
Type=forking
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{WORKSPACE_DIR}}
Environment=PATH={{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME={{HOME}}
EnvironmentFile={{INSTALL_DIR}}/config/config.env
EnvironmentFile=-{{CONFIG_DIR}}/config.env
ExecStart=/usr/bin/tmux -L lobster new-session -d -s lobster -c {{WORKSPACE_DIR}} {{CLAUDE_WRAPPER}}
ExecStop=/usr/bin/tmux -L lobster kill-session -t lobster
RemainAfterExit=yes
Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lobster-claude

[Install]
WantedBy=multi-user.target
