# services/lobster.target
[Unit]
Description=Lobster - Always-on Claude Code Message Processor
After=network.target
Wants=lobster-router.service amber-router.service lobster-claude.service

[Install]
WantedBy=multi-user.target
